# backend/migrations/systemd/luxquant-compass-resolver.service
[Unit]
Description=LuxQuant BTC Compass 2.0 — Projection First-Barrier Resolver
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/root/luxquant-terminal/backend
Environment="PYTHONUNBUFFERED=1"
ExecStart=/root/luxquant-terminal/backend/venv/bin/python -m app.services.compass_projection_resolver
StandardOutput=journal
StandardError=journal
TimeoutStartSec=300
