# ops/systemd/luxquant-btc-correlation-worker.service
[Unit]
Description=LuxQuant BTC Correlation Worker (LISTEN signal_created)
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/luxquant-terminal
EnvironmentFile=/root/luxquant-terminal/.env
ExecStart=/usr/bin/python3 /root/luxquant-terminal/btc_correlation_worker.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=luxquant-btc-correlation

[Install]
WantedBy=multi-user.target
