# deployment/luxquant-journey-worker.service
[Unit]
Description=LuxQuant Signal Journey Worker (real-time LISTEN mode)
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
WorkingDirectory=/root/luxquant-terminal/backend
ExecStart=/usr/bin/python3 -m app.workers.journey_worker
Restart=always
RestartSec=10
StandardOutput=append:/var/log/luxquant-sync/journey-worker.log
StandardError=append:/var/log/luxquant-sync/journey-worker.log

MemoryMax=2G

[Install]
WantedBy=multi-user.target
